# fixtures/launchd/org.rmxos.phase08.d22.keepalive-remove.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.rmxos.phase08.d22.keepalive-remove</string>
	<key>ProgramArguments</key>
	<array>
		<string>/root/nxplatform/phase1/d22-keepalive-remove-test</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
</dict>
</plist>
